# systemd/cpu-pl-set.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/cpu-pl-set

[Install]
WantedBy=multi-user.target
