# systemd/dns-update.timer
[Timer]
OnCalendar=*-*-* *:*:00
Unit=dns-update.service

[Install]
WantedBy=timers.target